# ThinImaging/ConfigureMunkiWebAdmin/Root/Library/LaunchDaemons/com.googlecode.munki.munkiwebadmin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>Label</key>
    <string>com.googlecode.munki.munkiwebadmin</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/Data/www/munkiwebadmin_env/bin/python</string>
        <string>/Volumes/Data/www/munkiwebadmin_env/munkiwebadmin/manage.py</string>
        <string>runwsgiserver</string>
        <string>port=8000</string>
        <string>host=0.0.0.0</string>
    </array>
    <key>StandardErrorPath</key>
    <string>/tmp/munkiwebadmin.log</string>
    <key>StandardOutPath</key>
    <string>/tmp/munkiwebadmin.log</string>
    <key>UserName</key>
    <string>osx</string>
  </dict>
</plist>